# brief/gunicornmccastellazzob.service
# ============================================
# Gunicorn Systemd Service per mccastellazzob.com
# ============================================
# Configurazione PRODUZIONE
# 
# Installazione:
#   sudo cp gunicornmccastellazzob.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable gunicornmccastellazzob
#   sudo systemctl start gunicornmccastellazzob
#
# Comandi utili:
#   sudo systemctl status gunicornmccastellazzob
#   sudo systemctl restart gunicornmccastellazzob
#   sudo journalctl -u gunicornmccastellazzob -f
# ============================================

[Unit]
Description=Gunicorn WSGI server per mccastellazzob.com
Documentation=https://docs.gunicorn.org/
After=network.target postgresql.service
Wants=postgresql.service

[Service]
# Utente e gruppo
User=gunicornweb2
Group=www-data

# Directory di lavoro
WorkingDirectory=/www/wwwroot/mccastellazzob.com/mccastellazzob

# Ambiente Python
Environment="PATH=/www/wwwroot/mccastellazzob.com/venv/bin:/usr/local/bin:/usr/bin"
Environment="DJANGO_SETTINGS_MODULE=mccastellazzob.settings.prod"
# Locale per supporto multilingua (it, en, fr)
Environment="LANG=it_IT.UTF-8"
Environment="LC_ALL=it_IT.UTF-8"
Environment="PYTHONIOENCODING=utf-8"

# Comando di avvio Gunicorn
# --workers: numero di worker (2 * CPU + 1 è la regola generale)
# --threads: thread per worker (utile per I/O bound)
# --timeout: timeout per richieste lunghe
# --max-requests: restart worker dopo N richieste (previene memory leak)
# --max-requests-jitter: variazione random per evitare restart simultanei
ExecStart=/www/wwwroot/mccastellazzob.com/venv/bin/gunicorn \
    --workers 3 \
    --threads 2 \
    --worker-class gthread \
    --bind unix:/www/wwwroot/mccastellazzob.com/mccastellazzob/mccastellazzob.sock \
    --timeout 120 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile /www/wwwlogs/mccastellazzob.com.gunicorn.access.log \
    --error-logfile /www/wwwlogs/mccastellazzob.com.gunicorn.error.log \
    --capture-output \
    --log-level info \
    mccastellazzob.wsgi:application

# Riavvio automatico in caso di crash
Restart=on-failure
RestartSec=5

# Limiti di sicurezza
# Previene scrittura su directory di sistema
ProtectSystem=strict
# Permette scrittura solo nelle directory specificate
ReadWritePaths=/www/wwwroot/mccastellazzob.com/mccastellazzob/media
ReadWritePaths=/www/wwwroot/mccastellazzob.com/mccastellazzob/cache
ReadWritePaths=/www/wwwlogs

# Limiti risorse (opzionali, decommentare se necessario)
# MemoryMax=512M
# CPUQuota=80%

# Permessi socket
UMask=0007

[Install]
WantedBy=multi-user.target
